# servers/gb.conf
[program:mongo_gb]
directory=servers/gb/services/mongo
command=bash -c "rm mongod.lock | echo && /opt/mongodb/bin/mongod --repair --dbpath ./ && exec /opt/mongodb/bin/mongod --dbpath ./ --port=27000"
process_name=%(program_name)s
stdout_logfile=servers/gb/services/logs/mongo_stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10 
stderr_logfile=servers/gb/services/logs/mongo_stderr.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
user=nobody
startsecs=10

[program:jsonic_gb]
directory=servers/gb/services/jsonic
command=python jsonic.py --port=80%(process_num)02d
process_name=%(program_name)s_%(process_num)02d
numprocs=2
stdout_logfile=servers/gb/services/logs/jsonic_stdout_%(process_num)02d.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=servers/gb/services/logs/jsonic_stderr_%(process_num)02d.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
first_port=8000
user=nobody

[program:torongo_gb]
directory=servers/gb/services/torongo
command=python jsonreststore.py --port=70%(process_num)02d --mongoport=27000 --logId=gb --logLevel=warning
process_name=%(program_name)s_%(process_num)02d
numprocs=2
stdout_logfile=servers/gb/services/logs/torongo_stdout_%(process_num)02d.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10 
stderr_logfile=servers/gb/services/logs/torongo_stderr_%(process_num)02d.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
first_port=7000
user=nobody
startretries=30

[program:miniServer_gb]
directory=servers/gb/services/miniServer
command=python server.py --port=71%(process_num)02d --logId=mini --logLevel=warning
process_name=%(program_name)s_%(process_num)02d
numprocs=1
stdout_logfile=servers/gb/services/logs/torongo_stdout_%(process_num)02d.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10 
stderr_logfile=servers/gb/services/logs/torongo_stderr_%(process_num)02d.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
first_port=7100
user=nobody

[group:gb]
programs=mongo_gb,jsonic_gb,torongo_gb,miniServer_gb

[nginx]
server_port=80
server_name=gb.cs.unc.edu gb localhost
http_insert: access_log /var/log/nginx/access.log vhcombined;
    error_log /var/log/nginx/error.log;
    gzip_static on;
server_insert:     location /playpen { alias /var/tmp/playpen; }
    location ~ /Media/Audio/([0-9a-f]+)([0-9a-f][0-9a-f])(.mp3|.ogg)$ {
        alias /var/Media/Audio/$2/$1$2$3;
    }
    location ~ /Media/Image/([0-9a-f]+)([0-9a-f][0-9a-f])(\..*)$ {
        alias /var/Media/Image/$2/$1$2$3;
    }

[component:mongo_gb]
init=mkdir services/mongo
program=mongo_gb
backup=servers/gb/services/mongo

[component:jsonic_gb]
init:git clone git://github.com/uncopenweb/jsonic.git components/jsonic;
symlinks:components/jsonic/server,services/jsonic, 
    components/jsonic/JSonic.js,www/libs/uow/audio/JSonic.js
program=jsonic_gb
proxies=/jsonic/,http://jsonic_gb/
aliases=/jsonic/files,services/jsonic/files

[component:torongo_gb]
init=ln -s /home/gb/Projects/NewWeb/torongo components/torongo
program=torongo_gb
symlinks=components/torongo/server,services/torongo, components/torongo/MongoStore.js,www/libs/uow/data/MongoStore.js
proxies=/data/,http://torongo_gb/data/
upload=/Media/upload/,/data/_upload/

[component:miniServer_gb]
init: git clone git://github.com/uncopenweb/miniServer.git components/miniServer
    mkdir -p /var/tmp/playpen
program=miniServer_gb
symlinks=components/miniServer,services/miniServer
proxies=/cgi/,http://miniServer_gb/

[component:dojotrace_gb]
init=git clone git://github.com/uncopenweb/dojotrace.git components/dojotrace
symlinks=components/dojotrace/trace.js,www/libs/uow/trace.js

[component:uow_gb]
init=git clone git://github.com/uncopenweb/uow.git components/uow; echo "gb" > www/libs/servername
symlinks=components/uow/uow.js,www/libs/uow.js, components/uow/uow/ui,www/libs/uow/ui, components/uow/catalog,www/catalog
aliases=/catalog,www/catalog

[component:local_gb]
aliases=/,/var/www/, /local,/home/gb/public_html
