# servers/hark.conf
[program:mongo_hark]
directory=servers/hark/services/mongo
command=bash -c "rm mongod.lock | echo && /opt/mongodb/bin/mongod --repair --dbpath ./ && exec /opt/mongodb/bin/mongod --dbpath ./ --port=27002"
process_name=%(program_name)s
stdout_logfile=servers/hark/services/logs/mongo_stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10 
stderr_logfile=servers/hark/services/logs/mongo_stderr.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
user=nobody
startsecs=10

[program:jsonic_hark]
directory=servers/hark/services/jsonic
command=python jsonic.py --port=%(process_num)d
process_name=%(program_name)s_%(process_num)d
numprocs=2
numprocs_start=8200
first_port=8200
stdout_logfile=servers/hark/services/logs/jsonic_stdout_%(process_num)d.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=servers/hark/services/logs/jsonic_stderr_%(process_num)d.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
user=nobody

[program:torongo_hark]
directory=servers/hark/services/torongo
command=python jsonreststore.py --port=%(process_num)d --mongoport=27002 --logId=hark --logLevel=debug
process_name=%(program_name)s_%(process_num)d
numprocs=2
numprocs_start=8210
first_port=8210
stdout_logfile=servers/hark/services/logs/torongo_stdout_%(process_num)d.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10 
stderr_logfile=servers/hark/services/logs/torongo_stderr_%(process_num)d.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
user=nobody
startretries=30

[group:hark]
programs=mongo_hark,jsonic_hark,torongo_hark

[nginx]
server_port=80
server_name=harkthesound.org harkthesound.com harkthesound.net local.hark
http_insert: access_log /var/log/nginx/access.log vhcombined;
    gzip_static on;
server_insert:    location ~ /Media/Audio/([0-9a-f]+)([0-9a-f][0-9a-f])(.mp3|.ogg)$ {
        alias /var/Media/Audio/$2/$1$2$3;
    }
    location ~ /Media/Image/([0-9a-f]+)([0-9a-f][0-9a-f])(\..*)$ {
        alias /var/Media/Image/$2/$1$2$3;
    }
    location ~ /nls/en(-us)?/ {
        echo "{}";
    }

[component:mongo_hark]
init=mkdir services/mongo
program=mongo_hark
backup=servers/hark/services/mongo

[component:jsonic_hark]
init=git clone git://github.com/uncopenweb/jsonic.git components/jsonic
symlinks=components/jsonic/server,services/jsonic, components/jsonic/JSonic.js,www/libs/uow/audio/JSonic.js
program=jsonic_hark
proxies=/jsonic/,http://jsonic_hark/
aliases=/jsonic/files,services/jsonic/files

[component:torongo_hark]
init:git clone git://github.com/uncopenweb/torongo.git components/torongo
    dd if=/dev/urandom of=components/torongo/server/secret bs=100 count=1
program=torongo_hark
symlinks=components/torongo/server,services/torongo, components/torongo/MongoStore.js,www/libs/uow/data/MongoStore.js
proxies=/data/,http://torongo_hark/data/
upload=/Media/upload/,/data/_upload/

[component:dojotrace_hark]
init=git clone git://github.com/uncopenweb/dojotrace.git components/dojotrace
symlinks=components/dojotrace/trace.js,www/libs/uow/trace.js

[component:uow_hark]
init=git clone git://github.com/uncopenweb/uow.git components/uow; echo "Hark the Sound" > www/libs/servername
symlinks=components/uow/uow.js,www/libs/uow.js, components/uow/uow/ui,www/libs/uow/ui, components/uow/catalog,www/catalog
aliases=/catalog,www/catalog

[component:hark_hark]
init=git clone git://github.com/uncopenweb/hark.org.git components/hark
symlinks=components/hark,www/hark
aliases=/,www/hark
